# systemd/server-boot-entry.path
[Unit]
Description=Watch for new ostree BLS entries and create server variant

[Path]
PathChanged=/boot/loader/entries

[Install]
WantedBy=paths.target
