# systemd/server-boot-entry-boot.service
[Unit]
Description=Ensure server boot entry exists at boot
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/create-server-boot-entry
